# conf/pccontrol.service
[Unit]
Description=Control volume from phone

[Service]
ExecStart=/home/vladimirov/go/bin/pccontrol_server

[Install]
WantedBy=multi-user.target

